# debian/sv/log/run
#!/bin/sh
exec setuidgid xplperllog multilog t s1000000 n1000 ./main
